# Rhizome/support/LaunchAgents/app.rhizome.watch-tick.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>app.rhizome.watch-tick</string>
    <!--
        AssociatedBundleIdentifiers is intentionally omitted: macOS 14+ enforces
        a launch-constraint check that requires the bundled program to be
        Team-ID-signed by the same team as the associated bundle. Adhoc-signed
        local builds fail that check with OS_REASON_CODESIGNING and never run.
        The cost is purely cosmetic — macOS attributes background activity to
        "compile-bin" instead of "Rhizome".
    -->
    <key>BundleProgram</key>
    <string>Contents/Resources/compile-bin</string>
    <key>ProgramArguments</key>
    <array>
        <string>compile-bin</string>
        <string>watch</string>
        <string>tick</string>
        <string>--json-stream</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
